# systemd/anatimer@cron.weekly.timer
[Unit]
Description=Run %E/cron.weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=multi-user.target
