# systemd/anatimer@.service
[Unit]
Description=Run scripts in directory '/etc/%i'

[Service]
Type=oneshot
ExecStart=/usr/local/bin/run-parts /etc/%i
