# systemd/anatimer@cron.daily.timer
[Unit]
Description=Run %E/cron.daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
