# systemd/anatimer@.service
[Unit]
Description=Run scripts in directory '%E/%i'

[Service]
Type=oneshot
ExecStartPre=mkdir -p '%E/%i'
ExecStart=run-parts '%E/%i'
