# systemd/anatimer@cron.monthly.timer
[Unit]
Description=Run %E/cron.monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=default.target
