# systemd/anatimer@cron.hourly.timer
[Unit]
Description=Run /etc/cron.hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=multi-user.target
